# .config/systemd/user/battery-warn.service
[Unit]
Description=Battery level warning check

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-warn
